# com.linkedin.leadagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.linkedin.leadagent</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>America/New_York</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>python3</string>
        <string>lead_agent.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tarunsrivastava/.gemini/antigravity/scratch/linkedin_pr_agent</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/tarunsrivastava/.gemini/antigravity/scratch/linkedin_pr_agent/lead_agent.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/tarunsrivastava/.gemini/antigravity/scratch/linkedin_pr_agent/lead_agent.err.log</string>
</dict>
</plist>
